# html_server.service
[Unit]
Description=html server

[Service]
ExecStart=/home/joe/secure_chat/html_server.py
