# http_server.service
[Unit]
Description=html server

[Service]
ExecStart=/home/joe/secure_chat/http_server.py
