# secure_chat_wss.service
[Unit]
Description=websocket server

[Service]
ExecStart=/home/joe/secure_chat/wss_server.py
