# ws_server.service
[Unit]
Description=websocket server

[Service]
ExecStart=/home/joe/secure_chat/ws_server.py
